# deployment/systemd/amilliontechies-backend.service.example
[Unit]
Description=Gunicorn instance for AmillionTechies backend
After=network.target

[Service]
User=dev
Group=dev
WorkingDirectory=/var/www/projects/amilliontechies/backend/backend
EnvironmentFile=/var/www/projects/amilliontechies/backend/.env
Environment=PATH=/var/www/projects/amilliontechies/backend/venv/bin
ExecStart=/var/www/projects/amilliontechies/backend/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8000 \
    backend.wsgi:application
Restart=always
RestartSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
